# gloom/srv/env.cron
1 * * * * export SCRIPT_PATH=/var/share/var/cron/main/ && cd $SCRIPT_PATH && ./init.sh > 00/cron.00.log
4 * * * * cd $CRONLOG_BASE_PATH && export FORECAST_PATH=/var/share/var/data/Channels/Forecast_Channel > 01/cron.01.log
4 * * * * cd $CRONLOG_BASE_PATH && export CMOC_PATH=/var/share/var/data/ > 02/cron.02.log
4 * * * * cd $CRONLOG_BASE_PATH && export EVC_PATH=/var/share/var/data/Channels/Everybody_Votes_Channel > 03/cron.03.log
4 * * * * cd $CRONLOG_BASE_PATH && export NEWS_PATH=/var/share/var/data/Channels/News_Channel > 04/cron.04.log
4 * * * * cd $CRONLOG_BASE_PATH && export NC_PATH=/var/share/var/data/Channels/Nintendo_Channel > 05/cron.05.log
9 * * * * cd $CRONLOG_BASE_PATH && export FLASK_APP=room.py > 06/cron.06.log
9 * * * * cd $CRONLOG_BASE_PATH && export FLASK_ENV=development > 07/cron.07.log
14 * * * * cd $BOT_PATH && ./auto.sh > 08/cron.08.log
